# {{cookiecutter.out_dir}}/sys/etc/cron.d/django.frep
PATH=/code/venv/bin:/code/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games
MAILTO=""

{{'{{'}} (((.Env.crons_mysupercron | default "false" )| lower ) | eq  "true" ) | ternary "" "#" -}}
0 7 5 * * django /bin/bash -c "/code/venv/bin/python /code/src/manage.py mysupercron &> /code/data/crons_mysupercron.log"
